# docs/openbsd/etc/rc.local
#/opt/solr-scripts/solr-start >/dev/null 2>&1 </dev/null # Start Solr
/etc/webmin/start >/dev/null 2>&1 </dev/null # Start Webmin
